# s6/s6-rc.d/redpanda/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid jobpod /usr/bin/redpanda \
  --redpanda-cfg /etc/redpanda/redpanda.yaml \
  --overprovisioned \
  --smp 1 \
  --reserve-memory 0M
