# s6/s6-rc.d/job-engine-cycle/run
#!/usr/bin/with-contenv bash
# Heartbeat loop: invoke `node dist/app/job.js cycle` every CYCLE_SECONDS.
# Each invocation is one-shot per spec §5.1 (fresh agent session per cycle).
# Non-zero exit propagates so s6-rc can mark the pod down (terminal state).
CYCLE_S="${CYCLE_SECONDS:-60}"
while true; do
  /usr/bin/node /opt/job_engine/dist/app/job.js cycle
  RET=$?
  if [[ $RET -ne 0 ]]; then
    echo "[cycle] node exited $RET" >&2
    exit $RET
  fi
  sleep "$CYCLE_S"
done
