# s6/s6-rc.d/job-engine-startup/up
/usr/bin/node /opt/job_engine/dist/app/job.js startup
