# s6/s6-rc.d/ingest-backfill/run
#!/usr/bin/with-contenv bash
cd /opt/up
exec s6-setuidgid jobpod /opt/venv/bin/python -m app.backfill
